# vuln_image/roles/fridge_service/files/fridge.service
[Unit]
Description=Fridge Service
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStart=/usr/bin/docker-compose -f /home/user/fridge/docker-compose.yml up
ExecStop=/usr/bin/docker-compose -f /home/user/fridge/docker-compose.yml down

[Install]
WantedBy=default.target
