# vuln_image/roles/doorlock_service/files/doorlock.service
[Unit]
Description=Doorlock Service
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStart=/usr/bin/docker-compose -f /home/user/doorlock/docker-compose.yml up
ExecStop=/usr/bin/docker-compose -f /home/user/doorlock/docker-compose.yml down

[Install]
WantedBy=default.target
