# vuln_image/roles/coffeepot_service/files/coffeepot.service
[Unit]
Description=Coffeepot Service
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStart=/usr/bin/docker-compose -f /home/user/coffeepot/docker-compose.yml up
ExecStop=/usr/bin/docker-compose -f /home/user/coffeepot/docker-compose.yml down

[Install]
WantedBy=default.target
